# deploy/ionos/systemd/ristoranti-site.service
[Unit]
Description=Ristoranti d'Italia site service
After=network.target

[Service]
Type=simple
User=ristoranti
Group=ristoranti
WorkingDirectory=/opt/ristoranti-site
EnvironmentFile=/etc/ristoranti/ristoranti.env
ExecStart=/usr/bin/node /opt/ristoranti-site/src/server/server.js
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ristoranti-site/data

[Install]
WantedBy=multi-user.target
