# deploy/ionos/systemd/ristoranti-bot-webhook.service
[Unit]
Description=Ristoranti d'Italia Telegram bot (webhook mode)
After=network.target

[Service]
Type=simple
User=ristoranti
Group=ristoranti
WorkingDirectory=/opt/ristoranti-site
EnvironmentFile=/etc/ristoranti/ristoranti.env
ExecStart=/usr/bin/node /opt/ristoranti-site/src/bot/bot.js
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
